# celery_app.conf
[program:celery_worker]
command=bash -c 'ENV_PATH=.env.prod celery -A sugar worker -l info'
directory=.
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=./logs/celery_worker.log

[program:celery_beat]
command=bash -c 'ENV_PATH=.env.prod celery -A sugar beat -l info'
directory=.
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=./logs/celery_beat.log

[supervisord]
logfile=./logs/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:12306

[inet_http_server]
port=127.0.0.1:12306
